# files/etc/systemd/system/connectivity-check.timer
[Unit]
Description=Run Daily Connectivity Check

[Timer]
OnBootSec=1d
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
