# files/etc/systemd/system/daily-connectivity-check.timer
[Unit]
Description=Run Daily Connectivity Check

[Timer]
OnBootSec=24h
OnUnitActiveSec=24h

[Install]
WantedBy=timers.target
