# files/etc/systemd/system/weekly-connectivity-check.timer
[Unit]
Description=Run Weekly Connectivity Check

[Timer]
OnBootSec=1w
OnUnitActiveSec=1w

[Install]
WantedBy=timers.target
