# etc/systemd/system/birdnet_data_migration.service
[Unit]
Description=BirdNET Data Migration Service
After=birdnet_recording.service
Requires=birdnet_recording.service

[Service]
ExecStart=/usr/bin/python3 /path/to/migration_script.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
